# supervisord.conf
[program:runserver]
exclude=true

[program: gunicorn]
directory = {{ APP_ROOT }}
command = python manage.py run_gunicorn
