# supervisord.conf
[program:runserver]
exclude=true

[program: gunicorn]
directory = {{ PROJECT_DIR }}
command = python manage.py run_gunicorn
