# supervisord.conf
[program:runserver]
exclude=true

[program: gunicorn]
directory = /mnt/xuflus/Webs/foowill
command = python manage.py run_gunicorn
